# book/experiments/hardened-runtime/out/launchctl/sandbox-lore.runtime-tools.ca072936-1257-4492-b1a7-89945d45a334.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONPATH</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime-tools.ca072936-1257-4492-b1a7-89945d45a334</string>
		<key>SANDBOX_LORE_CHANNEL</key>
		<string>launchd_clean</string>
		<key>SANDBOX_LORE_LAUNCHD_CLEAN</key>
		<string>1</string>
		<key>SANDBOX_LORE_RUN_ID</key>
		<string>ca072936-1257-4492-b1a7-89945d45a334</string>
		<key>SANDBOX_LORE_STAGE_OUTPUT_ROOT</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime-tools.ca072936-1257-4492-b1a7-89945d45a334/book/experiments/hardened-runtime/out</string>
		<key>SANDBOX_LORE_STAGE_ROOT</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime-tools.ca072936-1257-4492-b1a7-89945d45a334</string>
		<key>SANDBOX_LORE_STAGE_USED</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>sandbox-lore.runtime-tools.ca072936-1257-4492-b1a7-89945d45a334</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/opt/python@3.14/bin/python3.14</string>
		<string>-m</string>
		<string>book.api.runtime_tools</string>
		<string>run</string>
		<string>--plan</string>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime-tools.ca072936-1257-4492-b1a7-89945d45a334/book/experiments/hardened-runtime/plan.json</string>
		<string>--out</string>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime-tools.ca072936-1257-4492-b1a7-89945d45a334/book/experiments/hardened-runtime/out</string>
		<string>--channel</string>
		<string>direct</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime-tools.ca072936-1257-4492-b1a7-89945d45a334/book/experiments/hardened-runtime/out/launchctl/sandbox-lore.runtime-tools.ca072936-1257-4492-b1a7-89945d45a334.stderr.txt</string>
	<key>StandardOutPath</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime-tools.ca072936-1257-4492-b1a7-89945d45a334/book/experiments/hardened-runtime/out/launchctl/sandbox-lore.runtime-tools.ca072936-1257-4492-b1a7-89945d45a334.stdout.txt</string>
	<key>WorkingDirectory</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime-tools.ca072936-1257-4492-b1a7-89945d45a334</string>
</dict>
</plist>
