# book/experiments/runtime-closure/out/7edc2b2f-7edf-4a50-ba0c-bd9bb2a549d3/launchctl/sandbox-lore.runtime.db964c78-ca4b-467f-b2be-1748374f611d.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONPATH</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.db964c78-ca4b-467f-b2be-1748374f611d</string>
		<key>SANDBOX_LORE_CHANNEL</key>
		<string>launchd_clean</string>
		<key>SANDBOX_LORE_LAUNCHD_CLEAN</key>
		<string>1</string>
		<key>SANDBOX_LORE_RUN_ID</key>
		<string>db964c78-ca4b-467f-b2be-1748374f611d</string>
		<key>SANDBOX_LORE_SEATBELT_CALLOUT</key>
		<string>1</string>
		<key>SANDBOX_LORE_STAGE_OUTPUT_ROOT</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.db964c78-ca4b-467f-b2be-1748374f611d/book/experiments/runtime-closure/out/7edc2b2f-7edf-4a50-ba0c-bd9bb2a549d3</string>
		<key>SANDBOX_LORE_STAGE_ROOT</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.db964c78-ca4b-467f-b2be-1748374f611d</string>
		<key>SANDBOX_LORE_STAGE_USED</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>sandbox-lore.runtime.db964c78-ca4b-467f-b2be-1748374f611d</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/opt/python@3.14/bin/python3.14</string>
		<string>-m</string>
		<string>book.api.runtime</string>
		<string>run</string>
		<string>--plan</string>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.db964c78-ca4b-467f-b2be-1748374f611d/book/experiments/runtime-closure/plan.json</string>
		<string>--out</string>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.db964c78-ca4b-467f-b2be-1748374f611d/book/experiments/runtime-closure/out/7edc2b2f-7edf-4a50-ba0c-bd9bb2a549d3</string>
		<string>--channel</string>
		<string>direct</string>
		<string>--only-profile</string>
		<string>runtime-closure:iokit:v7_service_user_client_both</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.db964c78-ca4b-467f-b2be-1748374f611d/book/experiments/runtime-closure/out/7edc2b2f-7edf-4a50-ba0c-bd9bb2a549d3/launchctl/sandbox-lore.runtime.db964c78-ca4b-467f-b2be-1748374f611d.stderr.txt</string>
	<key>StandardOutPath</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.db964c78-ca4b-467f-b2be-1748374f611d/book/experiments/runtime-closure/out/7edc2b2f-7edf-4a50-ba0c-bd9bb2a549d3/launchctl/sandbox-lore.runtime.db964c78-ca4b-467f-b2be-1748374f611d.stdout.txt</string>
	<key>WorkingDirectory</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.db964c78-ca4b-467f-b2be-1748374f611d</string>
</dict>
</plist>
