# book/evidence/experiments/runtime-final-final/suites/runtime-checks/out/ca63d167-562c-4f3c-b332-6edf8683f9a8/launchctl/sandbox-lore.runtime.ca63d167-562c-4f3c-b332-6edf8683f9a8.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONPATH</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.ca63d167-562c-4f3c-b332-6edf8683f9a8</string>
		<key>SANDBOX_LORE_CHANNEL</key>
		<string>launchd_clean</string>
		<key>SANDBOX_LORE_LAUNCHD_CLEAN</key>
		<string>1</string>
		<key>SANDBOX_LORE_RUN_ID</key>
		<string>ca63d167-562c-4f3c-b332-6edf8683f9a8</string>
		<key>SANDBOX_LORE_STAGE_OUTPUT_ROOT</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.ca63d167-562c-4f3c-b332-6edf8683f9a8/book/evidence/experiments/runtime-final-final/suites/runtime-checks/out</string>
		<key>SANDBOX_LORE_STAGE_ROOT</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.ca63d167-562c-4f3c-b332-6edf8683f9a8</string>
		<key>SANDBOX_LORE_STAGE_USED</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>sandbox-lore.runtime.ca63d167-562c-4f3c-b332-6edf8683f9a8</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/opt/python@3.14/bin/python3.14</string>
		<string>-m</string>
		<string>book.api.runtime</string>
		<string>run</string>
		<string>--plan</string>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.ca63d167-562c-4f3c-b332-6edf8683f9a8/book/evidence/experiments/runtime-final-final/suites/runtime-checks/plan.json</string>
		<string>--out</string>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.ca63d167-562c-4f3c-b332-6edf8683f9a8/book/evidence/experiments/runtime-final-final/suites/runtime-checks/out</string>
		<string>--channel</string>
		<string>direct</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.ca63d167-562c-4f3c-b332-6edf8683f9a8/book/evidence/experiments/runtime-final-final/suites/runtime-checks/out/launchctl/sandbox-lore.runtime.ca63d167-562c-4f3c-b332-6edf8683f9a8.stderr.txt</string>
	<key>StandardOutPath</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.ca63d167-562c-4f3c-b332-6edf8683f9a8/book/evidence/experiments/runtime-final-final/suites/runtime-checks/out/launchctl/sandbox-lore.runtime.ca63d167-562c-4f3c-b332-6edf8683f9a8.stdout.txt</string>
	<key>WorkingDirectory</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.ca63d167-562c-4f3c-b332-6edf8683f9a8</string>
</dict>
</plist>
