# book/experiments/runtime-adversarial/out/922328dd-4ed3-45b3-be69-425b45b725f0/launchctl/sandbox-lore.runtime.922328dd-4ed3-45b3-be69-425b45b725f0.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONPATH</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.922328dd-4ed3-45b3-be69-425b45b725f0</string>
		<key>SANDBOX_LORE_CHANNEL</key>
		<string>launchd_clean</string>
		<key>SANDBOX_LORE_LAUNCHD_CLEAN</key>
		<string>1</string>
		<key>SANDBOX_LORE_RUN_ID</key>
		<string>922328dd-4ed3-45b3-be69-425b45b725f0</string>
		<key>SANDBOX_LORE_STAGE_OUTPUT_ROOT</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.922328dd-4ed3-45b3-be69-425b45b725f0/book/experiments/runtime-adversarial/out</string>
		<key>SANDBOX_LORE_STAGE_ROOT</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.922328dd-4ed3-45b3-be69-425b45b725f0</string>
		<key>SANDBOX_LORE_STAGE_USED</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>sandbox-lore.runtime.922328dd-4ed3-45b3-be69-425b45b725f0</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/opt/python@3.14/bin/python3.14</string>
		<string>-m</string>
		<string>book.api.runtime</string>
		<string>run</string>
		<string>--plan</string>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.922328dd-4ed3-45b3-be69-425b45b725f0/book/experiments/runtime-adversarial/plan.json</string>
		<string>--out</string>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.922328dd-4ed3-45b3-be69-425b45b725f0/book/experiments/runtime-adversarial/out</string>
		<string>--channel</string>
		<string>direct</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.922328dd-4ed3-45b3-be69-425b45b725f0/book/experiments/runtime-adversarial/out/launchctl/sandbox-lore.runtime.922328dd-4ed3-45b3-be69-425b45b725f0.stderr.txt</string>
	<key>StandardOutPath</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.922328dd-4ed3-45b3-be69-425b45b725f0/book/experiments/runtime-adversarial/out/launchctl/sandbox-lore.runtime.922328dd-4ed3-45b3-be69-425b45b725f0.stdout.txt</string>
	<key>WorkingDirectory</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.922328dd-4ed3-45b3-be69-425b45b725f0</string>
</dict>
</plist>
