# book/experiments/runtime-adversarial/out/launchctl/sandbox-lore.runtime-adversarial.clean.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONPATH</key>
		<string>/Users/achyland/Desktop/VM-BRIDGE/SANDBOX_LORE</string>
		<key>SANDBOX_LORE_SEATBELT_CALLOUT</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>sandbox-lore.runtime-adversarial.clean</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/python3</string>
		<string>/Users/achyland/Desktop/VM-BRIDGE/SANDBOX_LORE/book/experiments/runtime-adversarial/run_adversarial.py</string>
		<string>--require-clean</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/achyland/Desktop/VM-BRIDGE/SANDBOX_LORE/book/experiments/runtime-adversarial/out/launchctl/sandbox-lore.runtime-adversarial.clean.stderr.txt</string>
	<key>StandardOutPath</key>
	<string>/Users/achyland/Desktop/VM-BRIDGE/SANDBOX_LORE/book/experiments/runtime-adversarial/out/launchctl/sandbox-lore.runtime-adversarial.clean.stdout.txt</string>
	<key>WorkingDirectory</key>
	<string>/Users/achyland/Desktop/VM-BRIDGE/SANDBOX_LORE</string>
</dict>
</plist>
