# book/experiments/anchor-filter-map/out/launchctl/sandbox-lore.runtime-tools.4361513b-5bb0-40e4-aec7-5397b093c023.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONPATH</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime-tools.4361513b-5bb0-40e4-aec7-5397b093c023</string>
		<key>SANDBOX_LORE_CHANNEL</key>
		<string>launchd_clean</string>
		<key>SANDBOX_LORE_LAUNCHD_CLEAN</key>
		<string>1</string>
		<key>SANDBOX_LORE_RUN_ID</key>
		<string>4361513b-5bb0-40e4-aec7-5397b093c023</string>
		<key>SANDBOX_LORE_STAGE_OUTPUT_ROOT</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime-tools.4361513b-5bb0-40e4-aec7-5397b093c023/book/experiments/anchor-filter-map/out</string>
		<key>SANDBOX_LORE_STAGE_ROOT</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime-tools.4361513b-5bb0-40e4-aec7-5397b093c023</string>
		<key>SANDBOX_LORE_STAGE_USED</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>sandbox-lore.runtime-tools.4361513b-5bb0-40e4-aec7-5397b093c023</string>
	<key>ProgramArguments</key>
	<array>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime-tools.4361513b-5bb0-40e4-aec7-5397b093c023/.venv/bin/python</string>
		<string>-m</string>
		<string>book.api.runtime_tools</string>
		<string>run</string>
		<string>--plan</string>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime-tools.4361513b-5bb0-40e4-aec7-5397b093c023/book/experiments/anchor-filter-map/plan.json</string>
		<string>--out</string>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime-tools.4361513b-5bb0-40e4-aec7-5397b093c023/book/experiments/anchor-filter-map/out</string>
		<string>--channel</string>
		<string>direct</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime-tools.4361513b-5bb0-40e4-aec7-5397b093c023/book/experiments/anchor-filter-map/out/launchctl/sandbox-lore.runtime-tools.4361513b-5bb0-40e4-aec7-5397b093c023.stderr.txt</string>
	<key>StandardOutPath</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime-tools.4361513b-5bb0-40e4-aec7-5397b093c023/book/experiments/anchor-filter-map/out/launchctl/sandbox-lore.runtime-tools.4361513b-5bb0-40e4-aec7-5397b093c023.stdout.txt</string>
	<key>WorkingDirectory</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime-tools.4361513b-5bb0-40e4-aec7-5397b093c023</string>
</dict>
</plist>
