# book/experiments/runtime-closure/out/34a16e32-c4f9-4757-b69d-1b369a4d566e/launchctl/sandbox-lore.runtime.3bd540a9-ac52-4d3c-8321-6ba1d086be49.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONPATH</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.3bd540a9-ac52-4d3c-8321-6ba1d086be49</string>
		<key>SANDBOX_LORE_CHANNEL</key>
		<string>launchd_clean</string>
		<key>SANDBOX_LORE_LAUNCHD_CLEAN</key>
		<string>1</string>
		<key>SANDBOX_LORE_RUN_ID</key>
		<string>3bd540a9-ac52-4d3c-8321-6ba1d086be49</string>
		<key>SANDBOX_LORE_STAGE_OUTPUT_ROOT</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.3bd540a9-ac52-4d3c-8321-6ba1d086be49/book/experiments/runtime-closure/out/34a16e32-c4f9-4757-b69d-1b369a4d566e</string>
		<key>SANDBOX_LORE_STAGE_ROOT</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.3bd540a9-ac52-4d3c-8321-6ba1d086be49</string>
		<key>SANDBOX_LORE_STAGE_USED</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>sandbox-lore.runtime.3bd540a9-ac52-4d3c-8321-6ba1d086be49</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/opt/python@3.14/bin/python3.14</string>
		<string>-m</string>
		<string>book.api.runtime</string>
		<string>run</string>
		<string>--plan</string>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.3bd540a9-ac52-4d3c-8321-6ba1d086be49/book/experiments/runtime-closure/plan.json</string>
		<string>--out</string>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.3bd540a9-ac52-4d3c-8321-6ba1d086be49/book/experiments/runtime-closure/out/34a16e32-c4f9-4757-b69d-1b369a4d566e</string>
		<string>--channel</string>
		<string>direct</string>
		<string>--only-profile</string>
		<string>runtime-closure:iokit:v11_report_loud</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.3bd540a9-ac52-4d3c-8321-6ba1d086be49/book/experiments/runtime-closure/out/34a16e32-c4f9-4757-b69d-1b369a4d566e/launchctl/sandbox-lore.runtime.3bd540a9-ac52-4d3c-8321-6ba1d086be49.stderr.txt</string>
	<key>StandardOutPath</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.3bd540a9-ac52-4d3c-8321-6ba1d086be49/book/experiments/runtime-closure/out/34a16e32-c4f9-4757-b69d-1b369a4d566e/launchctl/sandbox-lore.runtime.3bd540a9-ac52-4d3c-8321-6ba1d086be49.stdout.txt</string>
	<key>WorkingDirectory</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.3bd540a9-ac52-4d3c-8321-6ba1d086be49</string>
</dict>
</plist>
