# book/experiments/runtime-closure/out/launchctl/sandbox-lore.runtime.a788dcc3-e484-4ae7-b764-f24bc86e5dfc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONPATH</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.a788dcc3-e484-4ae7-b764-f24bc86e5dfc</string>
		<key>SANDBOX_LORE_CHANNEL</key>
		<string>launchd_clean</string>
		<key>SANDBOX_LORE_IKIT_MACH_CAPTURE</key>
		<string>1</string>
		<key>SANDBOX_LORE_IKIT_METHOD0</key>
		<string>1</string>
		<key>SANDBOX_LORE_IKIT_METHOD0_PAYLOAD_IN</key>
		<string>book/experiments/runtime-closure/out/iosurface_method0_payload.plist</string>
		<key>SANDBOX_LORE_LAUNCHD_CLEAN</key>
		<string>1</string>
		<key>SANDBOX_LORE_RUN_ID</key>
		<string>a788dcc3-e484-4ae7-b764-f24bc86e5dfc</string>
		<key>SANDBOX_LORE_STAGE_OUTPUT_ROOT</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.a788dcc3-e484-4ae7-b764-f24bc86e5dfc/book/experiments/runtime-closure/out</string>
		<key>SANDBOX_LORE_STAGE_ROOT</key>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.a788dcc3-e484-4ae7-b764-f24bc86e5dfc</string>
		<key>SANDBOX_LORE_STAGE_USED</key>
		<string>1</string>
		<key>SBL_IKIT_SKIP_SWEEP</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>sandbox-lore.runtime.a788dcc3-e484-4ae7-b764-f24bc86e5dfc</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/opt/python@3.14/bin/python3.14</string>
		<string>-m</string>
		<string>book.api.runtime</string>
		<string>run</string>
		<string>--plan</string>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.a788dcc3-e484-4ae7-b764-f24bc86e5dfc/book/experiments/runtime-closure/plan.json</string>
		<string>--out</string>
		<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.a788dcc3-e484-4ae7-b764-f24bc86e5dfc/book/experiments/runtime-closure/out</string>
		<string>--channel</string>
		<string>direct</string>
		<string>--only-profile</string>
		<string>runtime-closure:iokit:v7_service_user_client_both</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.a788dcc3-e484-4ae7-b764-f24bc86e5dfc/book/experiments/runtime-closure/out/launchctl/sandbox-lore.runtime.a788dcc3-e484-4ae7-b764-f24bc86e5dfc.stderr.txt</string>
	<key>StandardOutPath</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.a788dcc3-e484-4ae7-b764-f24bc86e5dfc/book/experiments/runtime-closure/out/launchctl/sandbox-lore.runtime.a788dcc3-e484-4ae7-b764-f24bc86e5dfc.stdout.txt</string>
	<key>WorkingDirectory</key>
	<string>/private/tmp/sandbox-lore-launchctl/sandbox-lore.runtime.a788dcc3-e484-4ae7-b764-f24bc86e5dfc</string>
</dict>
</plist>
