# book/evidence/experiments/runtime-final-final/suites/mach-name-equivalence/second_opinion/3d9f98fd-e035-4b29-a2f0-5ce3d1c80def/lookup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>sandbox-lore.second-opinion.3d9f98fd-e035-4b29-a2f0-5ce3d1c80def.lookup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/private/tmp/sandbox_lore_second_opinion_3d9f98fd-e035-4b29-a2f0-5ce3d1c80def/mach_name_second_opinion</string>
    <string>lookup</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>WorkingDirectory</key><string>/private/tmp/sandbox_lore_second_opinion_3d9f98fd-e035-4b29-a2f0-5ce3d1c80def</string>
  <key>StandardOutPath</key><string>/private/tmp/sandbox_lore_second_opinion_3d9f98fd-e035-4b29-a2f0-5ce3d1c80def/lookup.stdout.txt</string>
  <key>StandardErrorPath</key><string>/private/tmp/sandbox_lore_second_opinion_3d9f98fd-e035-4b29-a2f0-5ce3d1c80def/lookup.stderr.txt</string>
</dict>
</plist>
