# book/evidence/experiments/runtime-final-final/suites/mach-name-equivalence/falsification/0b87ab83-6cff-4270-95c6-9f895f099096/job.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>sandbox-lore.falsification.0b87ab83-6cff-4270-95c6-9f895f099096</string>
  <key>ProgramArguments</key>
  <array>
    <string>/private/tmp/sandbox_lore_falsification_0b87ab83-6cff-4270-95c6-9f895f099096/sandbox_check_self_apply_matrix_probe</string>
    <string>/private/tmp/sandbox_lore_falsification_0b87ab83-6cff-4270-95c6-9f895f099096/profile.sb</string>
    <string>com.sandboxlore.falsify.globalname</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>WorkingDirectory</key><string>/private/tmp/sandbox_lore_falsification_0b87ab83-6cff-4270-95c6-9f895f099096</string>
  <key>StandardOutPath</key><string>/private/tmp/sandbox_lore_falsification_0b87ab83-6cff-4270-95c6-9f895f099096/stdout.txt</string>
  <key>StandardErrorPath</key><string>/private/tmp/sandbox_lore_falsification_0b87ab83-6cff-4270-95c6-9f895f099096/stderr.txt</string>
</dict>
</plist>
